# copy/opt/core/lib/svc/manifest/gogs.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest"
                name="export">
  <service name="application/gogs"
           type="service"
           version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network"
                grouping="require_all"
                restart_on="error"
                type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem"
                grouping="require_all"
                restart_on="error"
                type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <exec_method type="method"
                 name="start"
                 exec="/opt/core/lib/svc/method/gogs"
                 timeout_seconds="60">
      <method_context working_directory="/opt/gogs">
        <method_credential user="git" group="other" />
        <method_environment>
          <envvar name="HOME" value="/opt/gogs" />
          <envvar name="PATH" value="/bin:/usr/bin:/usr/local/bin:/opt/local/bin" />
          <envvar name="USER" value="git" />
        </method_environment>
      </method_context>
    </exec_method>
    <exec_method type="method"
                 name="stop"
                 exec=":kill"
                 timeout_seconds="60">
      <method_context/>
    </exec_method>
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Gogs: Go Git Service</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
